# deploy/com.engine.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the ENGINE scheduler daemon. Install:
    cp deploy/com.engine.scheduler.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.engine.scheduler.plist
  Adjust WorkingDirectory / paths to your checkout. RunAtLoad + KeepAlive keep it
  alive across logins; it detects wake and runs the morning chain (see scripts/scheduler.ts).
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.engine.scheduler</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>npx</string>
      <string>tsx</string>
      <string>scripts/scheduler.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yash/Desktop/Programming/fin_research</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/yash/Desktop/Programming/fin_research/data/logs/scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yash/Desktop/Programming/fin_research/data/logs/scheduler.err.log</string>
  </dict>
</plist>
